# crontab
# Refresh 4x/dia: 11:00 UTC (08:00 BRT), 17:00 UTC (14:00 BRT), 21:00 UTC (18:00 BRT), 23:00 UTC (20:00 BRT)
0 11,17,21,23 * * * /app/refresh.sh
